# tools/phantom-deployment-wizard/deployments/hidden-service/secure-deployment-user/systemd/deployment-watcher.path
[Unit]
Description=Watch for deployment packages
Documentation=man:systemd.path(5)

[Path]
PathChanged=/securepath/incoming/deployment_signed.zip
Unit=deployment-watcher.service

[Install]
WantedBy=multi-user.target