# phantom/scripts/phantom-multihop-interface.service
[Unit]
Description=Phantom Multihop VPN Interface Restore
After=network-online.target wg-quick@wg_main.service
Wants=network-online.target
ConditionPathExists=/opt/phantom-wg/config/phantom.json

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/bin/bash -c 'jq -e ".multihop.enabled == true" /opt/phantom-wg/config/phantom.json || exit 1'
ExecStart=/opt/phantom-wg/.phantom-venv/bin/python3 /opt/phantom-wg/phantom/scripts/multihop-interface-restore.py
ExecStop=/usr/bin/ip link del wg_vpn

# Security hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/phantom-wg /etc/wireguard
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW

# Resource limits
LimitNPROC=10
LimitNOFILE=100

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target