# phantom_frontmatter/systemd/phantom-frontmatter-ghost-wstunnel.service.template
[Unit]
Description=Phantom-WG Frontmatter Ghost wstunnel Server
Documentation=https://github.com/ARAS-Workspace/phantom-wg
After=network-online.target phantom-frontmatter-ghost-egress.service
Wants=network-online.target
Requires=phantom-frontmatter-ghost-egress.service
BindsTo=phantom-frontmatter-ghost-egress.service

[Service]
Type=simple
ExecStart=__BINARY_PATH__ server \
    wss://[::]:__LISTEN_PORT__ \
    __RESTRICT_PATH_ARG__ \
    --restrict-to __TARGET__ \
    --tls-certificate __CERT_PATH__ \
    --tls-private-key __KEY_PATH__
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=phantom-frontmatter-ghost-wstunnel
LimitNOFILE=65535

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
ReadWritePaths=__INSTALL_DIR__

[Install]
WantedBy=multi-user.target
